# upstart/kindle-trmnl.conf
# KindleTRMNL — Upstart auto-start job
#
# Install to: /etc/upstart/kindle-trmnl.conf
# (requires read-write root filesystem; already available on FW 5.16.x with hotfix)
#
# After copying, the job activates on next reboot automatically.
# To start immediately without rebooting:
#   start kindle-trmnl
# To stop:
#   stop kindle-trmnl

description "KindleTRMNL BYOS dashboard daemon"
author      "KindleTRMNL"

# Start after the Kindle GUI is fully up.
start on started lab126_gui
stop  on stopping lab126_gui

# Respawn up to 3 times within 30 seconds if the script crashes.
respawn
respawn limit 3 30

# Run as root (required for lipc, eips, and sysfs power control).
setuid root

script
    # Wait a few extra seconds for WiFi drivers to initialise.
    sleep 10
    exec /bin/sh /mnt/us/extensions/KindleTRMNL/trmnl.sh --autostart
end script

post-stop script
    # Ensure screen saver is re-enabled if the daemon dies unexpectedly.
    lipc-set-prop com.lab126.powerd preventScreenSaver 0 2>/dev/null || true
    lipc-set-prop com.lab126.cmd wirelessEnable 0 2>/dev/null || true
    rm -f /tmp/trmnl.lock /tmp/trmnl_refresh_flag
end script
